# system/service/sensorDriver.service
[Unit]
Description=Sensor Driver Service

[Service]
ExecStart=/home/guilherme/PycharmProjects/PJ2/system/driver.py
Restart=always